# home/supervisor.d/sc.conf
[program:sc]
umask=002
command=sc /home/conf/sc
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/tmp/sc.log
stdout_logfile_maxbytes=3MB
stdout_logfile_backups=2
user=root
